# infra/systemd/sentinel-watchdog.service
[Unit]
Description=Sentinel Watchdog
After=openclaw-sentinel.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sentinel-healthcheck
